# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:kooder]
autostart=true
autorestart=true
startretries=3
user=root
command=/root/bin/gateway.sh
